# simplified_system/deployment/supervisor/supervisord.conf
[supervisord]
nodaemon=true
user=quantuser
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid
minfds=65536

[unix_http_server]
file=/var/run/supervisor.sock

[supervisorctl]
serverurl=unix:///var/run/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[program:api_server]
command=python -m uvicorn src.api.main:app --host 0.0.0.0 --port 8000 --workers 8
directory=/app
user=quantuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/api_server.log
environment=PYTHONPATH=/app

[program:trading_engine]
command=python -m src.trading.main
directory=/app
user=quantuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/trading_engine.log
environment=PYTHONPATH=/app

[program:optimization_worker]
command=python -m src.optimization.worker
directory=/app
user=quantuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/optimization_worker.log
environment=PYTHONPATH=/app
numprocs=4

[program:monitoring_service]
command=python -m src.monitoring.prometheus_server
directory=/app
user=quantuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/monitoring_service.log
environment=PYTHONPATH=/app

[program:redis_cache]
command=redis-server --maxmemory 2gb --maxmemory-policy allkeys-lru
user=quantuser
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile=/app/logs/redis_cache.log